# deploy/oci/systemd/fusekit-worker-dispatch.service
[Unit]
Description=FuseKit hosted worker dispatch receiver
Documentation=https://github.com/Fusekit-AI/fusekit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fusekit
Group=fusekit
WorkingDirectory=/opt/fusekit
EnvironmentFile=/etc/fusekit/hosted-secrets.env
Environment=FUSEKIT_HOSTED_WORKER_ID=hosted-worker-dispatch
Environment=FUSEKIT_HOSTED_WORKER_WORKSPACE=/var/lib/fusekit/worker
Environment=FUSEKIT_HOSTED_WORKER_DISPATCH_STATE_DIR=/var/lib/fusekit/dispatch-state
ExecStart=/opt/fusekit/venv/bin/fusekit-hosted-worker-dispatch --host 127.0.0.1 --port 8766
Restart=on-failure
RestartSec=5s
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
StateDirectory=fusekit
StateDirectoryMode=0750
LogsDirectory=fusekit
LogsDirectoryMode=0750
RuntimeDirectory=fusekit
RuntimeDirectoryMode=0750
ReadWritePaths=/var/lib/fusekit /var/log/fusekit /run/fusekit

[Install]
WantedBy=multi-user.target
